# ansible/roles/tailscale/templates/tailscale/tailscaled.service.j2
# {{ ansible_managed }}

[Service]
LogLevelMax=notice
